# roles/gns/files/gns3.service
[Unit]
Description=GNS3 Server

[Service]
ExecStart=/usr/share/gns3/gns3-server/bin/gns3server

[Install]
WantedBy=multi-user.target
